# services/launchd/ai.openclaw.mesh-health-publisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.mesh-health-publisher</string>
  <key>ProgramArguments</key>
  <array>
    <string>${NODE_BIN}</string>
    <string>${HOME}/openclaw/bin/mesh-health-publisher.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCLAW_NODE_ID</key>
    <string>${OPENCLAW_NODE_ID}</string>
    <key>OPENCLAW_NATS</key>
    <string>${OPENCLAW_NATS}</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>${OPENCLAW_WORKSPACE}/.tmp/mesh-health-publisher.log</string>
  <key>StandardErrorPath</key>
  <string>${OPENCLAW_WORKSPACE}/.tmp/mesh-health-publisher.err</string>
</dict>
</plist>
