# services/launchd/ai.openclaw.log-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.log-rotate</string>
    <key>Comment</key>
    <string>Weekly log rotation for OpenClaw services</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/moltymac/openclaw/bin/log-rotate</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/moltymac/.openclaw/logs/log-rotate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moltymac/.openclaw/logs/log-rotate.log</string>
  </dict>
</plist>
