# services/launchd/ai.openclaw.mesh-task-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.mesh-task-daemon</string>
    <key>Comment</key>
    <string>OpenClaw Mesh Task Daemon — NATS task coordinator</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/moltymac/.openclaw/bin/node</string>
      <string>/Users/moltymac/openclaw/bin/mesh-task-daemon.js</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-task-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-task-daemon.err</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/moltymac</string>
      <key>PATH</key>
      <string>/Users/moltymac/.openclaw/workspace/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>NODE_PATH</key>
      <string>/Users/moltymac/openclaw/node_modules:/Users/moltymac/openclaw/lib</string>
      <key>OPENCLAW_NATS</key>
      <string>nats://100.91.131.61:4222</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
  </dict>
</plist>
