# services/launchd/ai.openclaw.mesh-deploy-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.mesh-deploy-listener</string>
  <key>ProgramArguments</key>
  <array>
    <string>${NODE_BIN}</string>
    <string>${HOME}/openclaw/bin/mesh-deploy-listener.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>${HOME}</string>
    <key>PATH</key>
    <string>${HOME}/.bun/bin:${HOME}/.local/bin:${HOME}/.npm-global/bin:${HOME}/bin:${HOME}/.volta/bin:${HOME}/.asdf/shims:${HOME}/Library/Application Support/fnm/aliases/default/bin:${HOME}/.fnm/aliases/default/bin:${HOME}/Library/pnpm:${HOME}/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OPENCLAW_NODE_ID</key>
    <string>${OPENCLAW_NODE_ID}</string>
    <key>OPENCLAW_NATS</key>
    <string>${OPENCLAW_NATS}</string>
    <key>OPENCLAW_NODE_ROLE</key>
    <string>${OPENCLAW_NODE_ROLE}</string>
    <key>OPENCLAW_REPO_DIR</key>
    <string>${OPENCLAW_REPO_DIR}</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>${OPENCLAW_WORKSPACE}/.tmp/mesh-deploy-listener.log</string>
  <key>StandardErrorPath</key>
  <string>${OPENCLAW_WORKSPACE}/.tmp/mesh-deploy-listener.err</string>
</dict>
</plist>
