# services/launchd/ai.openclaw.mission-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.mission-control</string>

    <key>WorkingDirectory</key>
    <string>/Users/moltymac/.openclaw/workspace/projects/mission-control</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/moltymac/.openclaw/bin/npm</string>
        <string>run</string>
        <string>dev</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>development</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mission-control.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mission-control.err</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
