# services/launchd/ai.openclaw.mission-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.mission-control</string>

    <key>WorkingDirectory</key>
    <string>${OPENCLAW_WORKSPACE}/projects/mission-control</string>

    <key>ProgramArguments</key>
    <array>
        <string>${HOME}/.openclaw/bin/npm</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>${HOME}</string>
        <key>PATH</key>
        <string>${HOME}/.bun/bin:${HOME}/.local/bin:${HOME}/.npm-global/bin:${HOME}/bin:${HOME}/.volta/bin:${HOME}/.asdf/shims:${HOME}/Library/Application Support/fnm/aliases/default/bin:${HOME}/.fnm/aliases/default/bin:${HOME}/Library/pnpm:${HOME}/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${OPENCLAW_WORKSPACE}/.tmp/mission-control.log</string>

    <key>StandardErrorPath</key>
    <string>${OPENCLAW_WORKSPACE}/.tmp/mission-control.err</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
