# services/launchd/ai.openclaw.mesh-deploy-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.mesh-deploy-listener</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/moltymac/.openclaw/bin/node</string>
    <string>/Users/moltymac/openclaw/bin/mesh-deploy-listener.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCLAW_NODE_ID</key>
    <string>moltymacs-virtual-machine-local</string>
    <key>OPENCLAW_NATS</key>
    <string>nats://100.91.131.61:4222</string>
    <key>OPENCLAW_NODE_ROLE</key>
    <string>lead</string>
    <key>OPENCLAW_REPO_DIR</key>
    <string>/Users/moltymac/openclaw-node</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-deploy-listener.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-deploy-listener.err</string>
</dict>
</plist>
