# services/launchd/ai.openclaw.mesh-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.mesh-agent</string>
    <key>Comment</key>
    <string>OpenClaw Mesh Agent Worker — claims and executes tasks via Claude Code CLI. NOT auto-loaded — load manually with: launchctl load ~/Library/LaunchAgents/ai.openclaw.mesh-agent.plist</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/moltymac/.openclaw/bin/node</string>
      <string>/Users/moltymac/openclaw/bin/mesh-agent.js</string>
    </array>
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moltymac/.openclaw/workspace/.tmp/mesh-agent.err</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/moltymac</string>
      <key>PATH</key>
      <string>/Users/moltymac/.openclaw/workspace/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>NODE_PATH</key>
      <string>/Users/moltymac/openclaw/node_modules:/Users/moltymac/openclaw/lib</string>
      <key>OPENCLAW_NATS</key>
      <string>nats://100.91.131.61:4222</string>
      <key>MESH_WORKSPACE</key>
      <string>/Users/moltymac/.openclaw/workspace</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
  </dict>
</plist>
