# services/systemd/openclaw-mesh-agent.service
[Unit]
Description=OpenClaw Mesh Agent — executes mesh tasks via the configured LLM provider (manually started)
After=network.target

[Service]
Type=simple
ExecStart=${NODE_BIN} ${OPENCLAW_REPO_DIR}/bin/mesh-agent.js
Restart=always
RestartSec=30
Environment=HOME=${HOME}
Environment=PATH=${HOME}/.openclaw/workspace/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_PATH=${OPENCLAW_REPO_DIR}/node_modules:${OPENCLAW_REPO_DIR}/lib
Environment=OPENCLAW_NATS=${OPENCLAW_NATS}
Environment=OPENCLAW_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=MESH_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=OPENCLAW_NODE_ID=${OPENCLAW_NODE_ID}
Environment=MESH_LLM_PROVIDER=${MESH_LLM_PROVIDER}
Environment=LLM_MODEL=${LLM_MODEL}
Environment=LLM_BASE_URL=${LLM_BASE_URL}
Environment=TZ=America/Montreal
WorkingDirectory=${OPENCLAW_REPO_DIR}
StandardOutput=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-agent.log
StandardError=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-agent.err

[Install]
WantedBy=default.target
