# services/systemd/openclaw-mesh-health-publisher.service
[Unit]
Description=OpenClaw Mesh Health Publisher — publishes node health metrics to NATS
After=network.target

[Service]
Type=simple
ExecStart=${NODE_BIN} ${OPENCLAW_REPO_DIR}/bin/mesh-health-publisher.js
Restart=always
RestartSec=5
Environment=HOME=${HOME}
Environment=PATH=${HOME}/.openclaw/workspace/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_PATH=${OPENCLAW_REPO_DIR}/node_modules:${OPENCLAW_REPO_DIR}/lib
Environment=OPENCLAW_NATS=${OPENCLAW_NATS}
Environment=OPENCLAW_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=OPENCLAW_NODE_ID=${OPENCLAW_NODE_ID}
Environment=TZ=America/Montreal
WorkingDirectory=${OPENCLAW_REPO_DIR}
StandardOutput=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-health-publisher.log
StandardError=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-health-publisher.err

[Install]
WantedBy=default.target
