# services/systemd/openclaw-scheduler-heartbeat.service
[Unit]
Description=OpenClaw Scheduler Heartbeat — POST /api/scheduler/tick so at-triggers and cron tasks fire without a Mission Control browser tab open
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s -o /dev/null --max-time 10 -X POST http://127.0.0.1:3000/api/scheduler/tick
