# services/systemd/openclaw-consolidation-scheduler.timer
[Unit]
Description=OpenClaw Consolidation Scheduler timer — every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
