# services/systemd/openclaw-lane-watchdog.service
[Unit]
Description=OpenClaw Lane Watchdog — monitors mesh lanes and restarts stalled services
After=network.target

[Service]
Type=simple
ExecStart=${NODE_BIN} ${HOME}/openclaw/bin/lane-watchdog.js
Restart=always
RestartSec=5
Environment=HOME=${HOME}
Environment=PATH=${HOME}/.openclaw/workspace/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_PATH=${HOME}/openclaw/node_modules:${HOME}/openclaw/lib
Environment=OPENCLAW_NATS=${OPENCLAW_NATS}
Environment=OPENCLAW_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=TZ=America/Montreal
WorkingDirectory=${HOME}/openclaw
StandardOutput=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-lane-watchdog.log
StandardError=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-lane-watchdog.err

[Install]
WantedBy=default.target
