# services/systemd/openclaw-mesh-tool-discord.service
[Unit]
Description=OpenClaw Mesh Tool Discord — Discord bot bridge for mesh commands
After=network.target

[Service]
Type=simple
ExecStart=${NODE_BIN} ${HOME}/openclaw/bin/mesh-tool-discord.js
Restart=always
RestartSec=5
Environment=HOME=${HOME}
Environment=PATH=${HOME}/.openclaw/workspace/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_PATH=${HOME}/openclaw/node_modules:${HOME}/openclaw/lib
Environment=OPENCLAW_NATS=${OPENCLAW_NATS}
Environment=OPENCLAW_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=TZ=America/Montreal
WorkingDirectory=${HOME}/openclaw
StandardOutput=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-tool-discord.log
StandardError=append:${OPENCLAW_WORKSPACE}/.tmp/openclaw-mesh-tool-discord.err

[Install]
WantedBy=default.target
