# services/systemd/openclaw-consolidation-scheduler.service
[Unit]
Description=OpenClaw Consolidation Scheduler — idle-checks then runs one memory consolidation cycle

[Service]
Type=oneshot
ExecStart=${NODE_BIN} ${OPENCLAW_WORKSPACE}/bin/consolidation-scheduler.mjs
Environment=HOME=${HOME}
Environment=PATH=${HOME}/.openclaw/workspace/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_WORKSPACE=${OPENCLAW_WORKSPACE}
Environment=OPENCLAW_NATS=${OPENCLAW_NATS}
Environment=OPENCLAW_NODE_ID=${OPENCLAW_NODE_ID}
Environment=TZ=America/Montreal
WorkingDirectory=${OPENCLAW_WORKSPACE}
StandardOutput=append:${OPENCLAW_WORKSPACE}/.tmp/consolidation-scheduler.log
StandardError=append:${OPENCLAW_WORKSPACE}/.tmp/consolidation-scheduler.err
