# services/nats/ai.openclaw.nats-1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.nats-1</string>
    <key>Comment</key>
    <string>OpenClaw NATS Node 1 — client 4222, cluster 6222, monitor 8222</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/nats-server</string>
      <string>--config</string>
      <string>${OPENCLAW_REPO}/services/nats/nats-1.conf</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/.openclaw/logs/nats-1.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.openclaw/logs/nats-1.err</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>${HOME}</string>
    </dict>
  </dict>
</plist>
